# deploy/attribut-generator.service
[Unit]
Description=Attribut Generator FastAPI Backend
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__REPO__/backend
Environment="PATH=__REPO__/.venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=__REPO__/backend/.env
ExecStart=__REPO__/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
UMask=0022
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__REPO__/backend /srv/images

[Install]
WantedBy=multi-user.target
